# .templates/supervisord/conf/workers/worker.conf
[program:worker]
command=php bin/console rabbitmq:consumer -w job
user=www-data
